# conf/supervisor/keepalived.ini
[program:keepalived]
command=keepalived -n -l         # -n: no daemonize, -l: log to syslog
autostart=true
autorestart=true
priority=2                       # lower priority than haproxy
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
